# bin/i3lock.service
[Unit]
Description=i3lock Service
Before=sleep.target

[Service]
Type=forking
User=mattshirlaw
Environment=DISPLAY=":1"
ExecStart=/usr/bin/i3lock -i /usr/share/backgrounds/warty-final-ubuntu.png

[Install]
WantedBy=suspend.target

